# board/poobs4/common/rootfs_overlay/etc/profile.d/env.sh
export CONFIGFS_HOME="/sys/kernel/config"
export PATH="/opt/keyboard-gadget:$PATH"
